# packages/rpi-agent/target/cumulocity-rpi-agent-8.19.0/control/postinst
chmod +x /usr/share/cumulocity-rpi-agent/*.sh

if [ -f "/etc/cumulocity-agent.properties" ]; then
    mv /usr/share/cumulocity-rpi-agent/cfg/cumulocity.properties /usr/share/cumulocity-rpi-agent/cfg/cumulocity.properties.template
else
    cp /usr/share/cumulocity-rpi-agent/cfg/cumulocity.properties /usr/share/cumulocity-rpi-agent/cfg/cumulocity.properties.template
    mv /usr/share/cumulocity-rpi-agent/cfg/cumulocity.properties /etc/cumulocity-agent.properties
fi
ln -s /etc/cumulocity-agent.properties /usr/share/cumulocity-rpi-agent/cfg/cumulocity.properties

if [ -x "/etc/init.d/cumulocity-agent" ]; then
    update-rc.d cumulocity-agent defaults 98 02 > /dev/null
    invoke-rc.d cumulocity-agent start || true
fi
